# ansible/roles/sprink/files/sprink.conf
[program:sprink]
command=/home/pi/git/sprink/run.py
autostart=true
autorestart=true
stderr_logfile=/var/log/sprink.err.log
stdout_logfile=/var/log/sprink.out.log